# files/default/ceilometer-agent-central-upstart.conf
description "Ceilometer Central Agent"
author "Jordan Tardif <jordan@dreamhost.com>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec /opt/ceilometer/bin/ceilometer-agent-central
